# com.shopify.evergreen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shopify.evergreen</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/samanthanorman/antigravity_skills/shopify-social-autoposter/main_facebook_instagram.py</string>
        <string>--evergreen</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Every Saturday at 10:00am -->
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/samanthanorman/antigravity_skills/shopify-social-autoposter/evergreen.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/samanthanorman/antigravity_skills/shopify-social-autoposter/evergreen_error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
