# com.shopify.folderwatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shopify.folderwatcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/samanthanorman/antigravity_skills/shopify-social-autoposter/folder_watcher.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/samanthanorman/antigravity_skills/shopify-social-autoposter</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Run at 9am, 12pm, and 4pm every day — posts up to 3 items total per day -->
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/samanthanorman/antigravity_skills/shopify-social-autoposter/folderwatcher_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/samanthanorman/antigravity_skills/shopify-social-autoposter/folderwatcher_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
